# Gunicorn.conf
[program:gunicorn]
command=gunicorn --workers=4 manage:app -b localhost:5000
directory=/deploy/app